# modules/GitHubRepoManager/GitHubRepoManager.TaskScheduler.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2024-08-01T04:00:00</Date>
    <Author>UnifiedAIToolbox</Author>
    <Description>Daily update and reporting for GitHub Repo Manager.</Description>
  </RegistrationInfo>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-18</UserId>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2024-08-01T04:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-NoProfile -ExecutionPolicy Bypass -Command "&amp; { Import-Module 'G:\Development\20_Staging\AI-Toolbox\UnifiedAIToolbox\modules\GitHubRepoManager\GitHubRepoManager.psm1'; Update-AllRepos; Export-RepoStatusReport }"</Arguments>
      <WorkingDirectory>G:\Development\20_Staging\AI-Toolbox\UnifiedAIToolbox</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
